# ubuntu1804_nginx_php_prod/root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv /bin/sh
exec /usr/local/php7/sbin/php-fpm -F
