# ubuntu1804_nginx_php_prod/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv /bin/sh
exec /usr/sbin/nginx -g 'daemon off;'
